# history/snapshots/20260303-215305/ai.openclaw.patch-check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.openclaw.patch-check</string>
    <key>Comment</key>
    <string>Re-apply OpenClaw patches if needed (after updates)</string>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>86400</integer>
    <key>WatchPaths</key>
    <array>
      <string>/opt/homebrew/lib/node_modules/openclaw/package.json</string>
    </array>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/mahsum/.openclaw/my-patches/post-update-check.sh</string>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/mahsum/.openclaw/logs/patch-check.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mahsum/.openclaw/logs/patch-check.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
      <key>HOME</key>
      <string>/Users/mahsum</string>
      <key>PATH</key>
      <string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
  </dict>
</plist>
